# service/chaekpool/ansible/roles/monitoring/files/grafana.openrc
#!/sbin/openrc-run

description="Grafana Dashboard Server"
command="/opt/grafana/bin/grafana-server"
command_args="--homepath=/opt/grafana --config=/etc/grafana/grafana.ini"
command_user="grafana:grafana"

supervisor="supervise-daemon"
respawn_max=0
respawn_delay=5
respawn_period=60
output_log="/var/log/grafana/grafana.log"
error_log="/var/log/grafana/grafana.log"

depend() {
    need net
    after prometheus loki jaeger
}
