# service/chaekpool/ansible/roles/valkey/files/valkey.openrc
#!/sbin/openrc-run

description="Valkey server"
command="/usr/bin/valkey-server"
command_args="/etc/valkey/valkey.conf"
command_user="valkey:valkey"

supervisor="supervise-daemon"
respawn_max=0
respawn_delay=5
respawn_period=60
output_log="/var/log/valkey/valkey.log"
error_log="/var/log/valkey/valkey.log"

depend() {
    need localmount net
}
