# service/chaekpool/scripts/valkey/configs/redis-commander.openrc
#!/sbin/openrc-run

description="Redis Commander Web Interface"
command="/usr/local/bin/redis-commander"
command_args="--redis-host 127.0.0.1 --redis-port 6379 --redis-password changeme --port 8081 --address 0.0.0.0"
command_user="rediscmd:rediscmd"

supervisor="supervise-daemon"
respawn_max=0
respawn_delay=5
respawn_period=60
output_log="/var/log/redis-commander.log"
error_log="/var/log/redis-commander.log"

depend() {
    need net
    after valkey
}
